# infra/systemd/credential-rotation.timer
[Unit]
Description=Weekly credential rotation for personal-stack

[Timer]
# Every Sunday at 03:00 UTC (after nightly pipeline at 02:00 UTC)
OnCalendar=Sun *-*-* 03:00:00 UTC
# Run missed rotations on next boot
Persistent=true
# Add up to 5 minutes of random jitter
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
