# infra/systemd/credential-rotation.service
[Unit]
Description=personal-stack credential rotation
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/personal-stack
ExecStart=/opt/personal-stack/infra/scripts/rotate-credentials.sh --all --force
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600

# Security hardening
ProtectSystem=strict
ReadWritePaths=/opt/personal-stack /var/run/docker.sock
PrivateTmp=true
NoNewPrivileges=false
