# infra/home-node/systemd/home-node-update.service
[Unit]
Description=Home node auto-update (git pull + re-configure)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/personal-stack/infra/home-node/update.sh
Environment=STACK_DIR=/opt/personal-stack
StandardOutput=journal
StandardError=journal
